# root/etc/services.d/dynomite/run
#!/usr/bin/with-contenv sh

exec /usr/local/bin/dynomite --conf-file=/etc/dynomite/dynomite.yml -v11;
